# udev/rules.d/50-firmware.rules
SUBSYSTEM=="firmware", ACTION=="add", RUN+="firmware.sh"
