# elements/dynamic-login/init-scripts/upstart/dynamic-login.conf
description     "Dynamic Login"

start on runlevel [2345]
stop on runlevel [!2345]

umask 022

expect stop

script
    echo "Executing Dynamic Login"
    /usr/local/bin/dynamic-login
end script
