# elements/dynamic-login/init-scripts/systemd/dynamic-login.service
[Unit]
Description=Dynamic Login
After=network.target

[Service]
ExecStart=/usr/local/bin/dynamic-login

[Install]
WantedBy=multi-user.target

