# elements/growroot/init-scripts/openrc/growroot
#!/sbin/runscript

start() {
    /usr/local/sbin/growroot
    eend 0
}

depend() {
    need localmount
}
