# elements/dhcp-all-interfaces/install.d/dhcp-interface@.service
[Unit]
Description=DHCP interface %I
Before=network-pre.target
Wants=network-pre.target

ConditionPathExists=!/etc/sysconfig/network-scripts/ifcfg-%I

[Service]
Type=oneshot
User=root
ExecStartPre=/usr/local/sbin/dhcp-all-interfaces.sh %I
ExecStart=/sbin/ifup %I
RemainAfterExit=true
TimeoutStartSec=30s

[Install]
WantedBy=multi-user.target
