# elements/growroot/init-scripts/upstart/growroot.conf
description "Growing root partition"

console output

start on local-filesystems

post-start script
/usr/local/sbin/growroot
end script
