# diskimage_builder/elements/ironic-agent/static/etc/systemd/system/ironic-agent-create-rescue-user.path
[Unit]
Description=Ironic user rescue - notify path existence

[Path]
PathExists=/etc/ipa-rescue-config/ipa-rescue-password

[Install]
WantedBy=multi-user.target
