# elements/mellanox/udev.d/81-mellanox-drivers.rules
ACTION!="add", GOTO="drivers_end"

SUBSYSTEM=="net", RUN+="/sbin/modprobe mlx4_en"

LABEL="drivers_end"

