# elements/config-applier/install.d/config-applier.upstart
start on runlevel [2345]

task

author "Clint Byrum <clint@fewbar.com>"
description "Config Applier"

exec config-applier
