# diskimage_builder/elements/growroot/init-scripts/openrc/growroot
#!/sbin/openrc-run

start() {
    /usr/local/sbin/growroot
    eend 0
}

depend() {
    need localmount
}
