# diskimage_builder/elements/runtime-ssh-host-keys/init-scripts/upstart/ssh-keygen.conf
description "OpenSSH Server Key Generation"

start on starting ssh
console output

task

exec /usr/local/sbin/runtime-ssh-host-keys.sh
